# create_admin_user.service
[Unit]
Description=Create CUPS admin user

[Service]
ExecStart=/opt/cups/run.sh

[Install]
WantedBy=multi-user.target
